# .config/systemd/user/aria2cd.service
[Unit]
Description=aria2 daemon

[Service]
ExecStart=/usr/bin/aria2c

[Install]
WantedBy=default.target
